# meta-quicksilver/recipes-quicksilver/simpleechoserver/simpleechoserver/simpleechoserver.service
[Unit]
Descriptioon=A simple echo server

[Service]
ExecStart=@BINDIR@/simple_echo_server 12345

[Install]
WantedBy=multi-user.target

